# docker/goaccess/supervisord.conf
[supervisord]
nodaemon=true

[program:cron]
command=service cron start

[program:webserver]
command=/start_webserver.sh